# test/extension/caching-proxy-parser/caching-proxy-parser.service
[Unit]
Description=Yuruna caching-proxy-parser (in-memory last-100 squid requests)
Documentation=https://yuruna.link/caching-proxy-parser
After=squid.service network-online.target
PartOf=squid.service

[Service]
Type=simple
User=proxy
Group=proxy
ExecStart=/usr/local/bin/caching-proxy-parser
Restart=on-failure
RestartSec=5
# proxy already owns /var/log/squid/yuruna_access.log (mode 0640) — the
# unit only needs to read it. No writes anywhere; the ring buffer is
# memory-only.
ReadOnlyPaths=/var/log/squid
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
